# innd-expire.service
[Unit]
Description=InterNet News Daemon Expiring

[Service]
User=news
PermissionsStartOnly=true
Type=simple
ExecStart=/usr/libexec/news/news.daily delayrm
